# site.list
/etc/doas.conf
/etc/exports
/etc/fstab.site
/etc/hostname.*
/etc/hosts.site
/etc/iked.conf
/etc/inetd.conf
/etc/installurl
/etc/mk.conf
/etc/mygate
/etc/myname
/etc/ntpd.conf
/etc/pf.conf
/etc/rc.conf.local
/etc/rc.d/nonxt_reflect
/etc/rc.firsttime
/etc/rc.local
/etc/relayd.conf
/etc/resolv.conf
/etc/snmpd.conf
/etc/soii.key
/etc/ssh/ssh_host_ecdsa_key
/etc/ssh/ssh_host_ecdsa_key.pub
/etc/ssh/ssh_host_ed25519_key
/etc/ssh/ssh_host_ed25519_key.pub
/etc/ssh/ssh_host_rsa_key
/etc/ssh/ssh_host_rsa_key.pub
/etc/ssh/ssh_known_hosts
/etc/sysctl.conf
/etc/vm.conf
/root/.cvsrc
/root/.profile
/root/.ssh/authorized_keys
/root/.ssh/id_rsa
/root/.ssh/id_rsa.pub
